# newhotspot.lodiwifi.net/etc/rc.local
# This executes on boot after most systems have been started.

# Start the directive processing loop
/root/directive.sh &

exit 0
